# deploy/systemd/bowery-agent.service
[Unit]
Description=The Bowery agent
Documentation=https://github.com/jvehent/the-bowery
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/bowery-agent --config /etc/bowery/agent.toml
Restart=on-failure
RestartSec=5s

# Resource isolation. Phase 0 ships a deliberately tight envelope; later phases
# raise these for the LLM analyzer (see DESIGN.md §6.4).
Slice=bowery.slice
MemoryHigh=512M
MemoryMax=1G
CPUWeight=100
TasksMax=256

# Capabilities. The eBPF loader needs CAP_BPF + CAP_PERFMON to load
# tracepoint programs and create ring-buffer maps. CAP_SYS_ADMIN is
# required as a fallback on kernels without CAP_BPF granularity
# (< 5.8). The agent drops to no privileges after BPF is loaded;
# Bowery never persists CAP_SYS_ADMIN in the run loop.
#
# Note: NoNewPrivileges interacts with eBPF — when set, only programs
# whose perf event creation is allowed for unprivileged users will
# load. For our tracepoints and ring buffers it's compatible. Keep
# the rest of the hardening intact.
CapabilityBoundingSet=CAP_BPF CAP_PERFMON CAP_SYS_ADMIN CAP_NET_ADMIN
AmbientCapabilities=CAP_BPF CAP_PERFMON CAP_SYS_ADMIN CAP_NET_ADMIN
NoNewPrivileges=yes

# Filesystem hardening.
ProtectSystem=strict
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
PrivateTmp=yes
PrivateDevices=yes
ReadWritePaths=/var/lib/bowery /var/log/bowery
ConfigurationDirectory=bowery
StateDirectory=bowery
LogsDirectory=bowery

# Network. Allow IPv4/IPv6 + Unix; later phases may tighten.
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX AF_NETLINK

# System call filter. Conservative default; Phase 2 may need to relax for eBPF.
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources

[Install]
WantedBy=multi-user.target
